# files/service/prometheus/run
#!/bin/sh
exec /opt/prometheus/prometheus \
  --config.file /etc/prometheus/prometheus.yml \
  --storage.tsdb.path /data \
  --web.listen-address 0.0.0.0:9090 \
  --storage.tsdb.retention 24h
