# files/service/grafana/run
#!/bin/sh
exec /opt/grafana/bin/grafana-server \
  -config /etc/grafana/config.ini \
  -homepath /opt/grafana
